# rc.local
sed -i "s@JENKINS_TOKEN@$JENKINS_TOKEN@g" /root/.kube/config
